# supervisord-foundationdb.conf
[supervisord]
nodaemon=true

[program:foundationdb]
command=/etc/init.d/foundationdb start
exitcodes=0
autorestart=false
startretries=0
startsecs=0
